# .config/autostart/brave-browser.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/brave-browser
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[ja_JP]=Brave Browser
Name=Brave Browser
Comment[ja_JP]=
Comment=
